# com.ike.ships-bell.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{USER}}.ships-bell</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>{{INSTALL_DIR}}/ships_bell.py</string>
        <string>--from</string>
        <string>{{START_HOUR}}</string>
        <string>--to</string>
        <string>{{END_HOUR}}</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>{{INSTALL_DIR}}</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>{{INSTALL_DIR}}/logs/ships-bell.log</string>
    
    <key>StandardErrorPath</key>
    <string>{{INSTALL_DIR}}/logs/ships-bell.error.log</string>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>LimitLoadToSessionType</key>
    <array>
        <string>Aqua</string>
    </array>
</dict>
</plist>